# cron.d/failsafe-watchdog
HOME=/
MAILTO=""

0 * * * * root /usr/local/bin/failsafe-watchdog
